# scripts/hrms.service
[Unit]
Description=HRMS Application Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/home/<your-user>/hrms-unified
ExecStart=/usr/bin/npm start
ExecStop=/usr/bin/npm run stop
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hrms-app
Environment=NODE_ENV=production
Environment=HOME=/home/<your-user>
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
