# scripts/hrms.service
[Unit]
Description=HRMS Application Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/home/<your-user>/hrms-unified
# Load environment variables including PORT from .env file
EnvironmentFile=/home/<your-user>/hrms-unified/.env
ExecStart=/usr/bin/npm run start:prod
ExecStop=/usr/bin/npm run stop
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hrms-app
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--openssl-legacy-provider
Environment=HOME=/home/<your-user>
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
